# examples/fw-restore.service
[Unit]
Description=Restore fw firewall rules at boot
Wants=network-online.target
After=network-online.target docker.service
ConditionPathExists=/path/to/firewall-manager/config/iptables.rules

[Service]
Type=oneshot
WorkingDirectory=/path/to/firewall-manager
ExecStart=/path/to/firewall-manager/fw restore
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
